# projects/控制台/launchd/com.yutu6.yuanxiao-decision-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.yutu6.yuanxiao-decision-pull</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/yutu6/.local/node-v24.16.0-darwin-arm64/bin/node</string>
    <string>/Users/yutu6/玉兔6工作区/projects/控制台/tools/pull-yuanxiao-decision.js</string>
  </array>
  <key>StartInterval</key><integer>90</integer>
  <key>StandardOutPath</key><string>/Users/yutu6/玉兔6工作区/logs/yuanxiao-decision-pull.log</string>
  <key>StandardErrorPath</key><string>/Users/yutu6/玉兔6工作区/logs/yuanxiao-decision-pull.log</string>
</dict>
</plist>
