# projects/控制台/launchd/com.yutu6.daily-governance-hardening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yutu6.daily-governance-hardening</string>

  <!-- launchd StartCalendarInterval 走本机本地时区。本机时区 = Asia/Shanghai(北京 UTC+8),
       故 Hour=5 即北京时间凌晨5点(= UTC 21:00)。若本机时区改为 UTC,需把 Hour 改成 21。 -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/yutu6/.local/node-v24.16.0-darwin-arm64/bin/node</string>
    <string>/Users/yutu6/玉兔6工作区/projects/控制台/tools/daily-governance-hardening.js</string>
    <string>--json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yutu6/玉兔6工作区/projects/控制台</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yutu6/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>YUTU6_NODE_BIN</key>
    <string>/Users/yutu6/.local/node-v24.16.0-darwin-arm64/bin/node</string>
    <key>CONSOLE_NODE_BIN</key>
    <string>/Users/yutu6/.local/node-v24.16.0-darwin-arm64/bin/node</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- 不在加载/开机时立即跑,只按日历准点触发,避免重复。 -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/yutu6/玉兔6工作区/projects/控制台/artifacts/daily-governance-hardening/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/yutu6/玉兔6工作区/projects/控制台/artifacts/daily-governance-hardening/launchd.err.log</string>
</dict>
</plist>
