# projects/控制台/launchd/com.yutu6.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yutu6.startup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/yutu6/玉兔6工作区/start-all.sh</string>
    <string>start</string>
    <string>--quiet</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/yutu6/玉兔6工作区</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/yutu6</string>
    <key>PATH</key>
    <string>/Users/yutu6/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>YUTU6_NODE_BIN</key>
    <string>/Users/yutu6/.local/node-v24.16.0-darwin-arm64/bin/node</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>LowPriorityIO</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/dev/null</string>
  <key>StandardErrorPath</key>
  <string>/Users/yutu6/玉兔6工作区/projects/控制台/artifacts/startup/launchd.err.log</string>
</dict>
</plist>
